# packaging/xdg-autostart/meridian-polkit-agent.desktop
[Desktop Entry]
Type=Application
Name=Meridian Polkit Agent
Comment=PolicyKit authentication agent for the Meridian desktop
Exec=/usr/local/bin/meridian-polkit-agent
OnlyShowIn=Meridian;
X-GNOME-Autostart-Phase=Initialization
